# sourcewell/crontab/newsletter.cron
59 23  *  *  * /usr/bin/wget http://sourcewell.berlios.de/nlsend.php3?period=daily -a /home/groups/sourcewell/htdocs/crontab/nldaily.log -O /home/groups/sourcewell/htdocs/crontab/nldaily.output
59 23  *  *  5 /usr/bin/wget http://sourcewell.berlios.de/nlsend.php3?period=weekly -a /home/groups/sourcewell/htdocs/crontab/nlweekly.log -O /home/groups/sourcewell/htdocs/crontab/nlweekly.output
